# root/target/linux/mediatek/base-files/etc/uci-defaults/99-net-ps
uci set network.globals.default_rps_val=14
uci set network.globals.default_rps_flow_cnt=256
uci set network.globals.default_xps_val=14
uci set network.globals.default_ps=1
uci commit
exit 0
