# root/target/linux/sunxi/base-files/etc/uci-defaults/99-switch-vlan-set-pvids
uci set network.@switch_vlan[0].vlan=10
uci set network.@switch_vlan[1].vlan=20
uci commit
exit 0
