# root/target/linux/sunxi/base-files/etc/uci-defaults/99-switch-config
if [ -f "/etc/.lamobo-r1.dsa" ]; then
	uci delete network.@switch_vlan[1]
	uci delete network.@switch_vlan[0]
	uci delete network.@switch[0]

	uci set network.lan.type=bridge
	uci set network.lan.ifname="wan"

	uci set network.wan1.type=bridge
	uci set network.wan2.type=bridge
	uci set network.wan3.type=bridge
	uci set network.wan4.type=bridge

	uci set network.wan1.ifname="lan1"
	uci set network.wan2.ifname="lan2"
	uci set network.wan3.ifname="lan3"
	uci set network.wan4.ifname="lan4"
else
	uci set network.@switch_vlan[0].vlan=10
	uci set network.@switch_vlan[1].vlan=20

	uci set network.@switch_vlan[0].ports="4 0 1 2 8t"
	uci set network.@switch_vlan[1].ports="3 8"
fi
uci commit
exit 0
