# root/target/linux/ipq40xx/base-files/etc/uci-defaults/05_add_mctl_cfg
#!/bin/sh

[ -n "$(uci -q get system.modem.disable)" ] || {
	for m in /sys/class/gpio/modem*_power; do
		local label="$(basename $m | awk -F_ '{print $1}')"
		uci set "system.${label}=mctl"
		uci set "system.${label}.disable=0"

		# modem is turned on in preinit but others are not
		[ "${label}" != "modem" ] && /sbin/mctl -p -m "${label}"
	done

	uci commit system
}

exit 0